# supervisord.conf
[supervisord]
nodaemon=true

[program:bot]
command=pipenv run bot

[program:web]
command=pipenv run web